# scripts/flyxiv-observer.service
[Unit]
Description=FlyXIV Observer Discord Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/flyxiv_observer
Environment=PATH=/home/ubuntu/flyxiv_observer/venv/bin
ExecStart=/home/ubuntu/flyxiv_observer/venv/bin/python -m pyobserver.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=flyxiv-observer

# Environment variables
Environment=NODE_ENV=production

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/ubuntu/flyxiv_observer

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target 